# calendar-service/assets/systemd-service-template.service
[Unit]
Description=Calendar Schedule Task

[Service]
ExecStart=@COMMAND@
Restart=on-failure
RestartSec=5s
